# .config/service/pipewire/run
#!/bin/sh
exec chpst -e "$TURNSTILE_ENV_DIR" pipewire 2>&1
