# .config/swww/service/swww-daemon/run
#!/bin/sh
exec chpst -e "$TURNSTILE_ENV_DIR" swww-daemon 2>&1