# .config/service/wireplumber/run
#!/bin/sh
exec chpst -e "$TURNSTILE_ENV_DIR" wireplumber 2>&1
