# .config/swww/service/swww-wallpaper/run
#!/bin/sh

exec chpst -e "$TURNSTILE_ENV_DIR" ./update-background 2>&1
