# .config/service/bluetooth-agent/run
#!/bin/sh

exec chpst -e "$TURNSTILE_ENV_DIR" bluetoothctl --timeout 60 agent on 2>&1