# .config/service/pipewire-pulse/run
#!/bin/sh
exec chpst -e "$TURNSTILE_ENV_DIR" pipewire-pulse 2>&1
